# pytel.service
[Unit]
Description=PyTel-Bot

[Service]
ExecStart=/home/pi/rpi_pytel_script.sh
User=pi
Group=pi

[Install]
WantedBy=network-online.target
