# pytel.service
[Unit]
Description=PyTel-Bot

[Service]
ExecStart=/home/pi/rpi_pytel_script.sh

[Install]
WantedBy=network-online.target
